# docs/userapisrv_test.ini
[program:userapisrv]
directory=/home/nizonglong/go/src/userapisrv/
command=/home/nizonglong/go/src/userapisrv/ -config=/home/nizonglong/go/src/userapisrv/test.toml
stdout_logfile=/data/log/userapisrv/userapisrv.log
stdout_logfile_backups=50
redirect_stderr=true
autostart=true
autorestart=true